# demo_systemd_gui/calc-forever/calc-forever.service
[Unit]
Description=galculator Forever
PartOf=graphical-session.target

[Service]
Restart=always
RestartSec=10
ExecStart=/usr/bin/galculator

[Install]
WantedBy=xsession.target
